# docker/bin/api/crontab
* * * * * root /usr/local/bin/php -q -f /app/artisan schedule:run --no-ansi >> /var/log/cron.log 2>&1
